# images/plasmon/import-pool.service
[Unit]
Description=Import zfs pool

[Service]
Type=oneshot
ExecStart=/usr/bin/zpool import pool

[Install]
WantedBy=multi-user.target
